# ansible/apps/jellyfin/templates/jellyfin-media-library.service.j2
[Unit]
Description=Jellyfin media library rclone mount
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
ExecStartPre=-/bin/fusermount -u {{ jellyfin_media_mount_path }}
ExecStart=/usr/bin/rclone mount "{{ jellyfin_media_remote }}" {{ jellyfin_media_mount_path }} --config {{ jellyfin_rclone_config_path }} --read-only --dir-cache-time=30m --poll-interval=30s --vfs-cache-mode=off
ExecStop=-/bin/fusermount -u {{ jellyfin_media_mount_path }}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
