# ansible/roles/jellyfin/templates/jellyfin.container.j2
[Unit]
Description=Jellyfin Media Server
After=network.target

[Container]
Image={{ jellyfin_image }}
PublishPort={{ jellyfin_port }}:8096
Volume={{ jellyfin_config_path }}:/config
Volume={{ jellyfin_cache_path }}:/cache
{% if jellyfin_media_dir is defined %}
Volume={{ jellyfin_media_dir }}:/media:ro
{% endif %}
UserNS=keep-id

[Service]
Restart=on-failure
# Rootless Podman storage init (storage-chown-by-maps) can take 10+ minutes
# on a Pi SD card the first time. RestartSec is long enough that a restart
# won't interrupt an in-progress chown.
RestartSec=60
TimeoutStartSec=1800

[Install]
WantedBy=default.target
